# debian/openct.postrm
#! /bin/sh

set -e

if [ "$1" = "purge" ]
then
    update-rc.d openct remove >/dev/null
fi

if [ "$1" = "remove" ]; then
        ldconfig
fi

exit 0
