# debian/openct.postinst
#! /bin/sh

set -e

if [ "$1" = "configure" ]
then
    update-rc.d openct defaults >/dev/null
    /etc/init.d/openct start
fi

if [ "$1" = "configure" ]; then
        ldconfig
fi

exit 0
